# services/chainlink_rewards_exporter.service
[Service]
Type=oneshot
User=root
ExecStart=/usr/local/bin/chainlink_rewards_exporter.sh
#TimeoutStartSec=3s
#TimeoutStopSec=1s
