# services/price_exporter.service
[Service]
Type=oneshot
User=root
ExecStart=/usr/local/bin/price_exporter.sh chainlink,ethereum
TimeoutStartSec=1s
TimeoutStopSec=1s
